# src/etc/crontabs/root
0 6 * * 0 /usr/share/v2ray/geo-update.sh