# src/etc/hotplug.d/iface/99-hiddify
#!/bin/sh

if [ "$INTERFACE" = "wan" ] || [ "$INTERFACE" = "wan6" ]; then
  if /etc/init.d/hiddify-cli status | grep -q "running"; then
    /etc/init.d/hiddify-cli restart
  fi
fi
