# src/etc/hotplug.d/iface/99-url-test
#!/bin/sh

[ "$ACTION" = "ifup" ] && /root/scripts/url-test.sh
