# src/etc/hotplug.d/iface/99-url-test
#!/bin/sh

[ "$ACTION" = "ifup" ] && \
[ "$DEVICE" = "$(ip route | grep default | awk -F 'dev ' '{print $2}' | awk '{print $1}')" ] && \
/root/scripts/url-test.sh
