# src/etc/hotplug.d/iface/99-warp
#!/bin/sh
[ "$INTERFACE" = "wan" ] || [ "$INTERFACE" = "wan6" ] || exit 0
/etc/init.d/warp-plus restart
/etc/init.d/warp-psiphon restart